# lj/docker/root/etc/supervisord.d/php.ini
[program:php]
command=/bin/bash -c "/usr/local/bin/composer install && php artisan migrate"